# src/services/ai.meta-factory.cortex.meta-factory.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.meta-factory.cortex.meta-factory</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/bin/cortex</string>
        <string>start</string>
        <string>--config</string>
        <string>__CONFIG_PATH__</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__CORTEX_DIR__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/.config/cortex/logs/cortex-meta-factory.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.config/cortex/logs/cortex-meta-factory.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__HOME__/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>GROVE_CHANNEL</key>
        <string>__AGENT_NAME__</string>
    </dict>
</dict>
</plist>
