# src/services/ai.meta-factory.cortex.relay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.meta-factory.cortex.relay</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BUN_PATH__</string>
        <string>__CORTEX_DIR__/src/taps/cc-events/relay.ts</string>
        <string>start</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__CORTEX_DIR__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/.claude/logs/cortex-relay.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.claude/logs/cortex-relay.error.log</string>
</dict>
</plist>
